# TW04/etc/crontab
#minute	hour	mday	month	wday	command
12	5	*	*	0	/bin/rm -f /var/log/*.html
